# docker/corp/mer-db02/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; Redis runs as root, no auth, all interfaces - intentionally exploitable.
[program:redis]
command=/usr/bin/redis-server --protected-mode no --bind 0.0.0.0 --port 6379 --dir /var/lib/redis
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:cacheadm]
command=/usr/bin/python3 /opt/cacheadm/cacheadm.py
user=cacheadm
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

; Seed the cache with loot keys a few seconds after redis comes up.
[program:redis-seed]
command=/bin/bash -c "sleep 5; redis-cli set ci:deploy_token 'mer-ci-9f3a17'; redis-cli set apikey:jenkins 'see mer-app01 (10.66.20.54) Jenkins for the deploy pipeline'; redis-cli set note 'cache for the CI pipeline - tokens reused on app01'"
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/dev/null
stderr_logfile=/dev/null
